# restic-prune.timer
[Unit]
Description=Prune data from the restic repository
[Timer]
OnCalendar=Mon *-*-* 9:30:00
Persistent=true
[Install]
WantedBy=timers.target
